# backend.service
[Unit]
Description=Gunicorn daemon for Django Standard Search Backend
After=network.target

[Service]
User=zkbz01
Group=zkbz01
WorkingDirectory=/home/zkbz01/std_searpage_v2/backend
# 注意：如果使用 Conda 虚拟环境，请将下方 ExecStart 路径替换为您的 Conda 环境路径，例如：/home/zkbz01/miniconda3/envs/std_searpage/bin/gunicorn
ExecStart=/home/zkbz01/std_searpage_v2/backend/venv/bin/gunicorn --workers 3 --bind 127.0.0.1:8001 config.wsgi:application
Restart=always
RestartSec=5
EnvironmentFile=/home/zkbz01/std_searpage_v2/backend/.env

[Install]
WantedBy=multi-user.target
